# tooling/systemd/mascope-disk-check.timer
[Unit]
Description=Run the Mascope disk-space monitor regularly

[Timer]
# Every 15 minutes: a runaway upload can fill a disk fast, and the check is
# cheap and read-only.
OnBootSec=5min
OnUnitActiveSec=15min
# Catch up one run if the machine was off at the scheduled time.
Persistent=true

[Install]
WantedBy=timers.target
