# tooling/systemd/mascope-update.service
[Unit]
Description=Mascope unattended update check
Documentation=https://github.com/karsa-oy/mascope
# An update only makes sense once Docker and the stack are up.
After=docker.service mascope.service
Requires=docker.service

[Service]
Type=oneshot
# The deploy user and mascope binary path below are filled in by
# tooling/ubuntu.sh, matching mascope.service. MASCOPE_PATH (and LD_PRELOAD)
# come from /etc/environment; update-specific settings and the release token
# live in the optional /etc/mascope/update.env (see docs/maintaining.md). The
# leading '-' means the unit still starts when that file is absent.
User=@@USER@@
EnvironmentFile=/etc/environment
EnvironmentFile=-/etc/mascope/update.env
ExecStart=@@MASCOPE_BIN@@ prod update --auto

# Exit 30 means "a migration update is pending" - informational, not a failure,
# so systemd must not mark the unit failed for it. 0 is a clean handled run.
SuccessExitStatus=0 30
